# platform/linux/systemd/claude-desktop-update.service
[Unit]
Description=Auto-update Claude Desktop from AUR
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/claude-desktop-update
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
