# platform/linux/systemd/claude-desktop-update.timer
[Unit]
Description=Check for Claude Desktop updates daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
